# scripts/kb-core.service
; systemd unit template for running the kb-core web UI on a Linux server.
; Edit the User, WorkingDirectory, and Environment lines, then install with:
;   sudo cp scripts/kb-core.service /etc/systemd/system/kb-core.service
;   sudo systemctl daemon-reload
;   sudo systemctl enable --now kb-core
;   sudo systemctl status kb-core
;   sudo journalctl -u kb-core -f       # live logs

[Unit]
Description=kb-core knowledge base web UI
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/pichia-kb
EnvironmentFile=/root/pichia-kb/.env
; ExecStart binds to 127.0.0.1 — access via SSH tunnel:
;   ssh -L 8501:localhost:8501 root@your-server
; To expose publicly (NOT recommended without nginx + auth in front),
; replace `127.0.0.1` with `0.0.0.0` and open port 8501 in security group.
ExecStart=/root/pichia-kb/.venv/bin/streamlit run web/Home.py \
    --server.port 8501 \
    --server.address 127.0.0.1 \
    --server.headless true \
    --browser.gatherUsageStats false
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
